# kcd-edinburgh-2025-10-22/demo_2/mnt.mount
[Mount]
What=demo
Where=/mnt
Type=virtiofs

[Install]
WantedBy=local-fs.target
